# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server --loglevel warning
autostart=true
autorestart=true
